# systemd-services/pdfjs.service
[Unit]
After=network.target

[Service]
Type=simple
ExecStart=/home/genshin/scripts/run_pdfjs.sh
User=genshin

[Install]
WantedBy=multi-user.target
